# mqttFan.service
[Unit]
Description=MQTT Fan
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3 /home/pi/Scripts/rpi-mqtt-fan/mqttFan.py

[Install]
WantedBy=multi-user.target
